# sys/supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:beanstalk]
command=/usr/bin/beanstalkd -f0 -b /data
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
startsecs=5

[program:beanstalkconsole]
command=php -S 0.0.0.0:8080 -t public
directory=/opt/beanstalk_console
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
startsecs=5
